# systemd/onlykey-agent.socket
# vim: set syn=systemd :
[Unit]
Description=onlykey-agent socket

[Socket]
ListenStream=/run/split-onlykey/%p/S.ssh-agent

ExecStartPost=-/bin/ln -s /run/split-onlykey/%p/S.ssh-agent /usr/local/etc/qubes-rpc/onlykey.SshAgent
ExecStopPost=-/bin/rm /usr/local/etc/qubes-rpc/onlykey.SshAgent
PassFileDescriptorsToExec=true
RemoveOnStop=true

[Install]
WantedBy=sockets.target


