# systemd/onlykey-agent.service
# vim: set syn=systemd :
[Unit]
Description=onlykey-agent

Requires=onlykey-agent.socket
Requires=onlykey-agent-daemon-factory.socket
After=onlykey-agent-daemon-factory.socket

[Service]
Type=notify
Environment=OK_SSH_SOCK=%t/%p/S.ssh-agent
ExecStart=/usr/local/bin/onlykey-agent-socket
RuntimeDirectory=split-onlykey/%p

# TODO: Sandboxing/hardening
