# systemd/onlykey-agent.socket
# vim: set syn=systemd :
[Unit]
Description=onlykey-agent socket

[Socket]
ListenStream=/run/split-onlykey/%p/ssh-agent.sock

Symlinks=/usr/local/etc/qubes-rpc/onlykey.SshAgent

[Install]
WantedBy=sockets.target


