# systemd/onlykey-agent-daemon-factory@.service
# vim: set syn=systemd :
[Unit]
Description=Onlykey SSH Agent daemon factory 

[Service]
ExecStart=onlykey-agent-daemon-factory
User=Root
